# monitoramento/linux/audit.rules
## -------------------------------------------------------------------------
## AUDIT RULES BASEADO EM NEO23X0 (Florian Roth) e socfortress
## https://github.com/Neo23x0/auditd/blob/master/audit.rules
## https://github.com/socfortress/Wazuh-Rules/blob/main/Auditd/auditd.conf
## -------------------------------------------------------------------------


# Remover regras existentes
-D

# Tamanho do Buffer (Aumentado para evitar perda em carga alta)
-b 8192

# Em caso de falha (0=silent, 1=printk, 2=panic)
-f 1

## r = read, w = write, a = mudança de atributo (permissões/metadados)

## -------------------------------------------------------------------------
## AUDIT RULES: MASTER CLASS PURPLE TEAM
## Fusão: Lab Custom + Florian Roth + Novas Sugestões
## -------------------------------------------------------------------------


## -------------------------------------------------------------------------
## 1. DETECÇÃO DE ATAQUES (ALTA PRIORIDADE - CASOS DE USO)
## -------------------------------------------------------------------------

# REVERSE SHELL (O "Pulo do Gato")
# Detecta se o bash tentar iniciar uma conexão de rede (saída para a internet)
-a always,exit -F arch=b64 -F exe=/bin/bash -F success=1 -S connect -k remote_shell
-a always,exit -F arch=b64 -F exe=/usr/bin/bash -F success=1 -S connect -k remote_shell
# (Versões 32 bits, caso necessário)
-a always,exit -F arch=b32 -F exe=/bin/bash -F success=1 -S connect -k remote_shell

# ESPIONAGEM (Credential Access)
-w /etc/shadow -p rwa -k shadow_access
-w /etc/gshadow -p rwa -k shadow_access
-w /etc/passwd -p wa -k identity_modification
-w /etc/group -p wa -k identity_modification

# DISCOVERY (Network Tools Execution)
-w /usr/bin/ncat -p x -k network_tool
-w /usr/bin/nc -p x -k network_tool
-w /bin/nc -p x -k network_tool
-w /usr/bin/netcat -p x -k network_tool
-w /usr/bin/nmap -p x -k network_tool
-w /usr/bin/socat -p x -k network_tool
-w /usr/bin/curl -p x -k network_tool
-w /usr/bin/wget -p x -k network_tool

## -------------------------------------------------------------------------
## 2. PERSISTÊNCIA & CONFIGURAÇÃO (Novas regras incluídas)
## -------------------------------------------------------------------------

# Shell Profiles (Onde malwares se escondem para iniciar com o usuário)
-w /etc/profile.d/ -p wa -k shell_profiles
-w /etc/profile -p wa -k shell_profiles
-w /etc/bashrc -p wa -k shell_profiles
-w /etc/csh.cshrc -p wa -k shell_profiles
-w /etc/zsh/ -p wa -k shell_profiles

# Network Modifications (Redirecionamento de DNS / Hosts)
-w /etc/hosts -p wa -k network_modifications
-w /etc/network/ -p wa -k network_modifications
-w /etc/issue -p wa -k network_modifications
-w /etc/hostname -p wa -k network_modifications
-a always,exit -F arch=b64 -S sethostname -S setdomainname -k network_modifications

# Cron (Persistência Clássica)
-w /etc/cron.d/ -p wa -k cron
-w /etc/cron.daily/ -p wa -k cron
-w /etc/cron.hourly/ -p wa -k cron
-w /etc/cron.weekly/ -p wa -k cron
-w /etc/cron.monthly/ -p wa -k cron
-w /etc/crontab -p wa -k cron
-w /var/spool/cron/crontabs/ -p wa -k cron

## -------------------------------------------------------------------------
## 3. RECONHECIMENTO (Discovery)
## -------------------------------------------------------------------------
-w /usr/bin/whoami -p x -k recon
-w /usr/bin/id -p x -k recon
-w /bin/hostname -p x -k recon
-w /bin/uname -p x -k recon

## -------------------------------------------------------------------------
## 4. SELF-DEFENSE
## Detecta tentativa de apagar logs ou parar a auditoria
## -------------------------------------------------------------------------
-w /var/log/audit/ -k audit_log_tampering
-w /etc/audit/ -p wa -k audit_config_change
-w /etc/libaudit.conf -p wa -k audit_config_change

## -------------------------------------------------------------------------
## 5. MONITORAMENTO GERAL (Base do Florian Roth)
## -------------------------------------------------------------------------
-a always,exit -F arch=b64 -S execve -k audit_exec
-a always,exit -F arch=b32 -S execve -k audit_exec
